# data/etc/supervisord.conf
; http://supervisord.org/configuration.html#file-format
; Open server on a port
; DEV: We use a port over UNIX socket due to lack of support by wercker
;   https://github.com/docker/docker/issues/12080
[inet_http_server]
port=127.0.0.1:2000

[supervisord]
; Set up our logs
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
; Keep pidfile same as `/etc/init.d`
pidfile=/var/run/supervisord.pid

[supervisorctl]
; Use our HTTP server as our server
serverurl=http://localhost:2000

; https://github.com/Supervisor/supervisor/blob/3.2.0/supervisor/skel/sample.conf#L41-L45
[rpcinterface:supervisor]
; Define protocol for communicating with supervisor server
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; Configure our services
; DEV: We use `-server` suffix to allow for multiple variations per service repo
;   For example, `-queue` and `-redis`
[program:app-server]
command=/usr/bin/node /home/ubuntu/app/main/bin/find-work-server
directory=/home/ubuntu/app/main
user=ubuntu
environment=ENV="production"

[program:app-queue]
command=/usr/bin/node /home/ubuntu/app/main/bin/find-work-queue
directory=/home/ubuntu/app/main
user=ubuntu
environment=ENV="production"

[program:app-redis]
command=/usr/bin/redis-server /etc/redis/app-redis.conf
user=ubuntu
